# nnnn.conf
[program:nnnn]
command=python3 /home/curator2/nnnn/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/nnnn/test.err.log
stdout_logfile=/home/curator2/nnnn/test.out.log